# systemd/basic.target
[Unit]
Description=Basic user session
After=sockets.target
Wants=sockets.target
