# systemd/graphical.target
[Unit]
Description=Graphical operation
After=basic.target
Wants=basic.target
