# systemd/headless.target
[Unit]
Description=Headless operation
After=basic.target
Wants=basic.target
