# systemd/threestatus.service
[Unit]
Description=threestatus daemon

[Service]
ExecStart=/home/janne/.config/i3/threestatus
