# systemd/threestatus.service
[Unit]
Description=threestatus daemon

[Service]
ExecStart=/home/janne/.dotfiles/graphical/threestatus
