# systemd/threestatus.socket
[Unit]
Description=threestatus socket

[Socket]
ListenStream=%t/threestatus.sock
SocketMode=0700

[Install]
WantedBy=sockets.target
